# docker/phpdocker/3party/config/supervisord.conf
[supervisord]
pidfile=/var/run/supervisord.pid
nodaemon=true
loglevel=debug
user=root

[program:nginx]
command=/usr/sbin/nginx -g daemon off
autostart=true
user=root
